# prod/uwsgi.conf
description "uwsgi instance"
start on runlevel [2345]
stop on runlevel [06]

respawn

exec /usr/bin/uwsgi --ini /etc/uwsgi.ini --logto /home/presentr/log.txt